# note-tagger_package/server_config/note_tagger_gun.conf
[program:gunicorn]
command=gunicorn note_tagger:application
directory: /home/ubuntu/note-tagger_package/note_tagger
user=nobody
autostart=true
autorestart=true
redirect_stderr=true